# telegram_story_bot/telegram-story-bot.service
[Unit]
Description=Telegram Story Bot
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/telegram_story_bot
ExecStart=/usr/bin/python3 /var/www/telegram_story_bot/telegram_story_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Переменные окружения
Environment=PYTHONUNBUFFERED=1

# Ограничения безопасности
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/telegram_story_bot

[Install]
WantedBy=multi-user.target
